# .config/systemd/user/daily_podman.timer
[Unit]
Description=Daily podman images updates

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
