# .config/systemd/user/daily_docker.timer
[Unit]
Description=Daily docker images updates

[Timer]
OnCalendar=1d

[Install]
WantedBy=timers.target
